# rpi-client/whisplay-agent.service
[Unit]
Description=Littlebat AI voice agent (Whisplay HAT)
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__SCRIPT_DIR__
EnvironmentFile=__SCRIPT_DIR__/config.env
ExecStart=/usr/bin/python3 __SCRIPT_DIR__/agent.py
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
